# fix-keyboard.service
[Unit]
Description=Fix RGB keyboard flicker on boot and resume
After=multi-user.target suspend.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/fix_keyboard.py

[Install]
WantedBy=multi-user.target suspend.target
